# conf/upstart/venise-cache-3.conf
author "Guilhem Saurel <gsaurel@laas.fr>"
chdir /home/gsaurel
respawn
respawn limit 10 60
setuid gsaurel
start on started network-services
stop on stopping network-services
env LANG="fr_FR.UTF-8"

description "Venise - Cache AGV 3"
exec python4 -m venise.sorties.cache -H yuki
