# conf/upstart/venise-cache-1.conf
author "Guilhem Saurel <gsaurel@laas.fr>"
chdir /home/gsaurel
respawn
respawn limit 10 60
setuid gsaurel
start on started venise-trajectoire-granier
stop on stopping venise-trajectoire-granier
env LANG="fr_FR.UTF-8"

description "Venise - Cache AGV 1"
exec python3 -m venise.sorties.cache -H moro
