# conf/upstart/venise-log-ubi-3.conf
author "Guilhem Saurel <gsaurel@laas.fr>"
chdir /home/gsaurel
respawn
respawn limit 10 60
setuid gsaurel
start on started networking
stop on stopping networking
env LANG="fr_FR.UTF-8"

description "Venise - Logging AGV positions"
exec python3 -m venise.sorties.log_ubisense -H yuki
