# conf/upstart/venise-granier-interieur.conf
author "Guilhem Saurel <gsaurel@laas.fr>"
chdir /home/gsaurel
respawn
respawn limit 10 60
setuid gsaurel
start on started network-services
stop on stopping network-services


description "Venise - Sonde granier intérieure"
exec python3 -m venise.sorties.granier_interieur
